# systemd-service-file/game-server.service
[Unit]
Description=Game Server

[Service]
Type=simple
ExecStart=/root/server --daemon
Restart=always

[Install]
WantedBy=multi-user.target
